# sandbox/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command= mongod --fork --dbpath /var/lib/mongodb/ --smallfiles --logpath /var/log/mongodb.log --logappend
stderr_logfile=/dev/stderr

